# s6/service-iplayer-arr/run
#!/usr/bin/env bash
exec s6-setuidgid hotio /app/iplayer-arr
